# .config/systemd/user/clock.timer
[Unit]
Description=Handle eww widgets

[Timer]
OnCalendar=*-*-* *:*:*
AccuracySec=1us
Unit= clock.service

[Install]
WantedBy=timers.target

