# .config/systemd/user/clock.service
[Unit]
Description=Handle eww widgets

[Service]
Type=oneshot
ExecStart=%h/.config/eww/scripts/open_timedate.sh

