# .config/systemd/user/theme_changer.service
[Unit]
Description=Check theme at boot

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/theme_change/check_time.sh

